# deploy/orchid-tracker.service
[Unit]
Description=OrchidTracker - Orchid Collection Manager
After=network.target
Wants=network.target

[Service]
Type=simple
User=orchid
Group=orchid
WorkingDirectory=/opt/orchids
ExecStart=/opt/orchids/target/release/orchid-tracker
EnvironmentFile=/opt/orchids/.env
Environment=LEPTOS_TAILWIND_VERSION=v4.2.0
Environment=RUST_LOG=info
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orchid-tracker

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/orchids/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
